# agent/installer/darwin/payload/Library/LaunchDaemons/com.openlabstats.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openlabstats.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/openlabstats/openlabstats-agent</string>
        <string>--config</string>
        <string>/usr/local/openlabstats/configs/agent.yaml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/var/log/openlabstats/agent-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/openlabstats/agent-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/openlabstats</string>
</dict>
</plist>
